# board/raspberrypi64-mass-storage-gadget/rootfs_overlay/etc/systemd/system/serial-getty@ttyS0.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --keep-baud --autologin root 921600,460800,230400,115200,57600,38400,9600 - $TERM
